# sv/keybase/run
#!/bin/sh
exec 2>&1
exec keybase service --chdir ~/.config/keybase
