# sv/keybase/run
#!/bin/sh
exec 2>&1
exec keybase --debug service --chdir ~/.config/keybase
