# sv/vbox-win/run
#!/bin/mksh

exec VBoxSDL --startvm "Windows 8.1"
