# sv/vbox-win/run
#!/bin/sh

exec VBoxSDL --startvm "Windows 8.1"
