# sv/unclutter/run
#!/usr/bin/sh

exec unclutter --jitter 10 --timeout 3
