# sv/unclutter/run
#!/usr/bin/sh

exec unclutter -noevents -root -visible -jitter 10 -idle 5
