# sv/unclutter/run
#!/usr/bin/mksh

exec unclutter -noevents -root -visible -jitter 10 -idle 5
