# sv/pipewire-pulse/run
#!/bin/sh
sv check ~/service/dbus-session >/dev/null 2>&1 || exit 1
exec pipewire-pulse
