# sv/statusbar/run
#!/bin/sh

exec 2>&1
exec /home/duncan/go/bin/go-dwmstatus
