# sv/statusbar/run
#!/bin/mksh

exec ~/repos/go-dwmstatus/go-dwmstatus
