# sv/dbus-session/run
#!/bin/execlineb -P
importas -i XDG_RUNTIME_DIR XDG_RUNTIME_DIR
dbus-daemon
 --session
 --nofork
 --nopidfile
 --address=unix:path=${XDG_RUNTIME_DIR}/bus
