# sv/pipewire-filters/run
#!/bin/sh
exec pipewire -c filter-chain.conf
