# sv/wallpaper/run
#!/bin/sh

feh --recursive --randomize --bg-fill --no-fehbg ~/.wallpaper
exec chpst -b wallpaper pause
