# sv/snooze-fstrim/run
#!/bin/sh

exec 2>&1
exec snooze -n -t ./timefile fstrim -av
