# sv/runsvdir-duncan/finish
#!/bin/sh

sv -w600 force-stop /home/duncan/service/*
sv exit /home/duncan/service/*
