# sv/runsvdir-duncan/run
#!/bin/sh

[ ! -d /run/runit.duncan ] && \
  mkdir -p /run/runit.duncan && \
  chown duncan:duncan /run/runit.duncan

export HOME=/home/duncan
groups=$(id -Gn duncan | tr ' ' ':')

exec 2>&1
exec chpst -u "duncan:$groups" sh -l -c "exec runsvdir /home/duncan/service 'log: ...........................................................................................................................................................................................................................................................................................................................................................................................................'" duncan
