# sv/runsvdir-duncan/run
#!/bin/sh

exec 2>&1
exec chpst -e ./env -u duncan:duncan runsvdir /home/duncan/service/
