# sv/ssh-agent/finish
#!/bin/sh

rm -f "${XDG_CACHE_HOME:-$HOME/.cache}"/ssh-agent.sock
exit 0
