# sv/ssh-agent/run
#!/bin/sh

exec ssh-agent -Da "${XDG_CACHE_HOME:-$HOME/.cache}"/ssh-agent.sock
