# sv/snooze-xlocate/run
#!/bin/sh

exec 2>&1
exec snooze -H13 -s 1d -t ./timefile sh -c 'xlocate -S ; touch ./timefile'
