# sv/ck-session/run
#!/bin/mksh

[ -d "./env" ] && exit 1

mkdir ./env

save_env() {
  while IFS= read -r var; do
    echo "${var#*=}" >> "./env/${var%%=*}"
  done
}

/usr/bin/ck-launch-session env \
  | grep "XDG_SESSION_COOKIE\|XDG_RUNTIME_DIR" \
  | save_env

exec pause
