# sv/x/finish
#!/bin/sh

sv -w600 force-stop ~/service.x/*
sv exit ~/service.x/*
rm -rf /tmp/serverauth.*
